# packaging/jobview.cron
MAILTO=
*/2 * * * * apache jobview-update
